# optimal_scheduler/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare Consumer_asset_IDs
declare Generator_asset_IDs
declare Consumer_asset_IDs

## Get the the user config options.
Consumer_asset_IDs=$(bashio::config 'Consumer_asset_IDs')
Generator_asset_IDs=$(bashio::config 'Generator_asset_IDs')
Energy_Source_asset_IDs=$(bashio::config 'Energy_Source_asset_IDs')
Building_consumption_IDs=$(bashio::config 'Building_consumption_IDs')

## Run your program
python Abstraction/main.py "$SUPERVISOR_TOKEN" "$Consumer_asset_IDs" "$Generator_asset_IDs" "$Energy_Source_asset_IDs" "$Building_consumption_IDs"

exit