# optimal_scheduler/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare message
declare toggle_something

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')
toggle_something=$(bashio::config 'Toggle_something')
Consumer_asset_IDs=$(bashio::config 'Consumer asset IDs')

## Run your program
python Abstraction/main.py "$SUPERVISOR_TOKEN" "$message" "$Consumer_asset_IDs"

exit